# layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.benny.mylogd.plist
